# taskwarrior/.config/systemd/user/tasksync.timer
[Unit]
Description=Taskwarrior Sync timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=15m
Unit=tasksync.service

[Install]
WantedBy=timers.target
